# templates/update-konstructs.socket.j2
[Unit]
Description=Remote Konstructs Server update trigger

[Socket]
ListenStream=44711
Accept=yes

[Install]
WantedBy=sockets.target
